# config/init/mailcatcher.conf
description "Mailcatcher"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/local/rvm/bin/mailcatcher --foreground --http-ip=0.0.0.0
